# kpcl-automation.service
[Unit]
Description=KPCL Automation System
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/kpcl-automation
Environment=PATH=/home/ec2-user/kpcl-automation/venv/bin
ExecStart=/home/ec2-user/kpcl-automation/venv/bin/python app.py
Restart=always
RestartSec=10

# Environment variables
EnvironmentFile=/home/ec2-user/kpcl-automation/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ec2-user/kpcl-automation

[Install]
WantedBy=multi-user.target
